# bind-formula/bind/files/bind-container.service.jinja
# bind-container.service, generated by bind formula
{%- set key_directory = salt['pillar.get']('bind:config:key_directory', map.key_directory) %}

[Unit]
Description=Bind container service
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Environment=BIND_IMAGE="{{ map.image }}"
Restart=on-failure
ExecStartPre=/bin/rm -f %t/bind-container.pid %t/bind-container.ctr-id

ExecStart=/bin/sh -c '/usr/bin/podman run \
        -p 53:53/udp -p 53:53/tcp \
        -v /etc/named.d:/etc/named.d:ro,Z \
        -v /etc/named.conf:/etc/named.conf:ro,Z \
        -v {{ map.named_directory }}:{{ map.container_named_directory }}:Z \
        -v {{ key_directory }}:{{ key_directory }}:ro,Z \
	--conmon-pidfile %t/bind-container.pid \
	--cidfile %t/bind-container.ctr-id \
	--cgroups=no-conmon \
	--replace -dt \
	--name bind-container \
	${BIND_IMAGE} '

ExecStop=/usr/bin/podman stop --ignore --cidfile %t/bind-container.ctr-id -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/bind-container.ctr-id
PIDFile=%t/bind-container.pid
TimeoutStopSec=60
Type=forking

[Install]
WantedBy=multi-user.target default.target
