# dhcpd-formula/dhcpd/files/kea-container.service.jinja
# kea-container.service, generated by dhcp formula

[Unit]
Description=Kea container service
Wants=network.target
After=network-online.target

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Environment=KEA_IMAGE="{{ dhcpd.image }}"
Restart=on-failure
ExecStartPre=/bin/rm -f %t/kea-container.pid %t/kea-container.ctr-id

ExecStart=/bin/sh -c '/usr/bin/podman run \
        --network=host --cap-add=NET_RAW \
        -v /etc/kea:/etc/kea:ro,Z \
	--conmon-pidfile %t/kea-container.pid \
	--cidfile %t/kea-container.ctr-id \
	--cgroups=no-conmon \
	--replace -dt \
	--name kea-container \
	${KEA_IMAGE} \
	kea-dhcp4 -c /etc/kea/kea-dhcp4.conf'

ExecStop=/usr/bin/podman stop --ignore --cidfile %t/kea-container.ctr-id -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/kea-container.ctr-id
PIDFile=%t/kea-container.pid
TimeoutStopSec=60
Type=forking

[Install]
WantedBy=multi-user.target default.target
